# s6/unicorn/finish
rm /voc/tmp/pids/unicorn.pid

/bin/true
