# s6/unicorn/run
#!/usr/bin/with-contenv /bin/sh

exec bundle exec unicorn -c config/unicorn.rb
