# s6/rake-reporting-nightly/run
#!/usr/bin/with-contenv /bin/sh

cd /voc

exec bundle exec rake resque:work
