# s6/rake-reporting-nightly/run
#!/usr/bin/with-contenv /bin/sh

sleep 30s
cd /voc

exec bundle exec rake reporting:daily
